# backend/app/config/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn --config config/gunicorn_config.py backend:app

# [program:flask]
# command=python backend.py