# frontend/supervisord.conf
[supervisord]
nodaemon=true

[program:streamlit]
command=python APP_wrap.py